# agent-daemon.service
[Unit]
Description=Claude Agent Daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Set AGENT_ID to this machine's identifier
Environment=AGENT_ID=thinkpad
Environment=MCP_URL=https://mcp.howling.one
Environment=DAEMON_PORT=9100
# Optional: set WEBHOOK_HOST to Tailscale IP if auto-detection fails
# Environment=WEBHOOK_HOST=100.x.x.x
ExecStart=/usr/bin/python3 /opt/agent-daemon/agent_daemon.py
WorkingDirectory=/opt/agent-daemon
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
